# third_party/jitterentropy/jitterentropy-library/tests/raw-entropy/recording_restart_kernelspace/boottime_test_record.service
# Systemd configuration file
#
[Unit]
Description=Boot time test for Kernel Jitter RNG
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target

[Service]
ExecStart=/usr/local/sbin/boottime_test_record.sh

[Install]
WantedBy=basic.target
